# docker-compose/pattern-3/supervisor/supervisord-gateway-manager.conf
[supervisord]
nodaemon=true

[program:sshd]
user=root
startsecs = 120
autorestart = false
command=/usr/sbin/sshd -D
